# orm-rules-tests/slow_server@.service
[Unit]
Description=Slow HTTP server on port %I
After=network.target

[Service]
Type=simple
ExecStart=/root/slow_server.py --port %i --name %i

[Install]
WantedBy=multi-user.target
